# deploy/systemd/jarvis-madusa-cycle.service
[Unit]
Description=MADUSA daily discovery cycle
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
User=jarvis
WorkingDirectory=/opt/jarvis
ExecStart=/opt/jarvis/jarvis madusa cycle
